# daemon/scripts/com.maiocchi.zeusdaemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maiocchi.zeusdaemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__ZEUSDAEMON_BINARY__</string>
        <string>--port</string>
        <string>2223</string>
        <string>--host</string>
        <string>127.0.0.1</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>/tmp/zeusdaemon.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/zeusdaemon.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
